# root/etc/s6-overlay/s6-rc.d/decap/run
#!/command/execlineb -P
exec s6-setuidgid abc sh -c 'cd /app && yarn serve'