# root/etc/s6-overlay/s6-rc.d/config/run
#!/command/execlineb -P
echo "No config necessary"